# recording/mqtt-recorder-reconcile.timer
[Unit]
Description=Jadwal reconcile CSV -> Supabase tiap 6 jam

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
